# cron/photozoom-analytics.cron
0 */2 * * * "/Users/aduvarov/Documents/PhotoZoom Analytics/scripts/run_daily_report.sh" >> "/Users/aduvarov/Documents/PhotoZoom Analytics/report.log" 2>&1
